# scripts/teslacam.service
[Unit]
Description=TeslaCam
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/teslacam

[Install]
WantedBy=multi-user.target
